# backlight.rules
RUN+="/bin/chmod a+w /sys/devices/platform/hp-wmi/rgb_zones/zone00"
RUN+="/bin/chmod a+w /sys/devices/platform/hp-wmi/rgb_zones/zone01"
RUN+="/bin/chmod a+w /sys/devices/platform/hp-wmi/rgb_zones/zone02"
RUN+="/bin/chmod a+w /sys/devices/platform/hp-wmi/rgb_zones/zone03"
